# eclipse-apps/etc/eclipse/systemd/system/drm-init.service
[Unit]
Description=DRM Display Initialization Service
After=basic.target

[Service]
Type=oneshot
ExecStart=/usr/bin/drm-init
RemainAfterExit=yes

[Install]
WantedBy=hardware.target
