# eclipse-apps/etc/eclipse/systemd/system/eclipse-shell.service
[Unit]
Description=Eclipse Shell Service
After=basic.target network.service

[Service]
Type=simple
ExecStart=/usr/bin/eclipse-shell
Restart=on-failure

[Install]
WantedBy=multi-user.target
