# eclipse-apps/s6/services/syslog/run
#!/bin/sh
# Syslog Service - S6 run script
# Replaces: syslog.service

exec 2>&1
exec /sbin/eclipse-syslog
