# eclipse-apps/s6/services/network/run
#!/bin/sh
# Network Service - S6 run script
# Replaces: network.service

exec 2>&1
exec /sbin/eclipse-network
