# eclipse-apps/s6/services/eclipse-shell/run
#!/bin/sh
# Eclipse Shell Service - S6 run script
# Replaces: eclipse-shell.service

exec 2>&1
exec /sbin/eclipse-shell
